# server-conf/supervisor-nginx.conf
[program:nginx]
command = service nginx start
stdout_events_enabled=true
stderr_events_enabled=true
environment=LANG="en_US.UTF-8",LC_ALL="en_US.UTF-8",LC_LANG="en_US.UTF-8"

